# overlay/etc/services.d/freshclam/run
#!/bin/bash

echo " - running clamav updater..."
freshclam -d -c 6
