# overlay/etc/services.d/clamd/run
#!/bin/bash

echo " - running clamav daemon..."
clamd -c /etc/clamav/clam.conf
