# overlay/etc/services.d/clamd/run
#!/bin/sh

#======================================================================================================================
# Run clamav in daemon mode
#   -c = use specified configuration file
#======================================================================================================================

echo " - running clamav daemon..."
clamd -c /etc/clamav/clamd.conf
