# overlay/etc/services.d/clamd/run
#!/bin/sh

echo " - running clamav daemon..."
clamd -c /etc/clamav/clam.conf
