# overlay/etc/services.d/clamd/run
#!/usr/bin/with-contenv sh

set -euo pipefail
export BF_E="${PWD##*/}/$(basename ${0})"


#======================================================================================================================
# Run clamav in daemon mode.
#
# Clamd flags
#   -c  use specified configuration file
#======================================================================================================================

bf-echo "Starting ClamAV daemon."
clamd -c ${CLAMAV_CLAMD_CONF}
