# overlay/etc/services.d/clamd/run
#!/usr/bin/with-contenv sh

#======================================================================================================================
# Run clamav in daemon mode
#   -c = use specified configuration file
#======================================================================================================================

_echo "starting clamav daemon"
clamd -c /etc/clamav/clamd.conf
