# overlay/etc/s6-overlay/s6-rc.d/clamd/run
#!/bin/withenv sh

set -euo pipefail


#======================================================================================================================
# Run preflight checks
#======================================================================================================================

nu -c "use bf-clamav run ; run preflight"


#======================================================================================================================
# Start ClamAV in foreground mode
#======================================================================================================================

clamd --foreground --config-file=${BF_CLAMAV_CLAMD_CONF}
